# root/etc/s6-overlay/s6-rc.d/init-ssh-config/run
#!/usr/bin/with-contenv bash
# SSH and Claude Code initialization script
# Handles persistence of SSH keys, user password, and Claude credentials

# ============================================
# Username Configuration
# ============================================
USER_NAME="${USER_NAME:-abc}"

# Rename user if custom username is specified
if [ "$USER_NAME" != "abc" ]; then
    echo "[init-ssh-config] Renaming user 'abc' to '$USER_NAME'..."
    usermod -l "$USER_NAME" abc
    groupmod -n "$USER_NAME" abc
    usermod -d "/home/$USER_NAME" -m "$USER_NAME"
fi

# Create persistent directories
mkdir -p /config/ssh
mkdir -p /config/claude-profile
mkdir -p /config/ssh/user_keys

# ============================================
# SSH Host Key Persistence
# ============================================
# Check if persistent host keys exist
if [ ! -f /config/ssh/ssh_host_ed25519_key ]; then
    echo "[init-ssh-config] Generating new SSH host keys..."
    ssh-keygen -t ed25519 -f /config/ssh/ssh_host_ed25519_key -N ""
    ssh-keygen -t rsa -b 4096 -f /config/ssh/ssh_host_rsa_key -N ""
    chown root:root /config/ssh/ssh_host_*
    chmod 600 /config/ssh/ssh_host_*_key
    chmod 644 /config/ssh/ssh_host_*_key.pub
else
    echo "[init-ssh-config] Using existing SSH host keys from /config/ssh/"
fi

# Remove any existing ephemeral host keys and symlink persistent ones
rm -f /etc/ssh/ssh_host_*
ln -sf /config/ssh/ssh_host_ed25519_key /etc/ssh/ssh_host_ed25519_key
ln -sf /config/ssh/ssh_host_ed25519_key.pub /etc/ssh/ssh_host_ed25519_key.pub
ln -sf /config/ssh/ssh_host_rsa_key /etc/ssh/ssh_host_rsa_key
ln -sf /config/ssh/ssh_host_rsa_key.pub /etc/ssh/ssh_host_rsa_key.pub

# ============================================
# User SSH Keys Persistence
# ============================================
# Setup authorized_keys persistence for user
if [ ! -d "/home/$USER_NAME" ]; then
    mkdir -p "/home/$USER_NAME"
fi

if [ -d "/home/$USER_NAME/.ssh" ] && [ ! -L "/home/$USER_NAME/.ssh" ]; then
    # If .ssh exists and is not a symlink, move contents to persistent storage
    cp -rn "/home/$USER_NAME/.ssh/"* /config/ssh/user_keys/ 2>/dev/null || true
    rm -rf "/home/$USER_NAME/.ssh"
fi

ln -sf /config/ssh/user_keys "/home/$USER_NAME/.ssh"
chown -R "$USER_NAME:$USER_NAME" /config/ssh/user_keys 2>/dev/null || true

# ============================================
# Claude Code Profile Persistence
# ============================================
# Symlink Claude config directory for persistence
if [ -d "/home/$USER_NAME/.claude" ] && [ ! -L "/home/$USER_NAME/.claude" ]; then
    cp -rn "/home/$USER_NAME/.claude/"* /config/claude-profile/ 2>/dev/null || true
    rm -rf "/home/$USER_NAME/.claude"
fi

ln -sf /config/claude-profile "/home/$USER_NAME/.claude"
chown -R "$USER_NAME:$USER_NAME" /config/claude-profile

# ============================================
# Password Configuration
# ============================================
if [ -n "$PASSWORD" ]; then
    echo "[init-ssh-config] Setting password for user '$USER_NAME' from PASSWORD environment variable"
    echo "$USER_NAME:$PASSWORD" | chpasswd
else
    # Generate random password if not provided
    RANDOM_PASS=$(tr -dc 'A-Za-z0-9' < /dev/urandom | head -c 16)
    echo "$USER_NAME:$RANDOM_PASS" | chpasswd
    echo "============================================"
    echo "[init-ssh-config] WARNING: No PASSWORD env var set!"
    echo "[init-ssh-config] Generated random password for '$USER_NAME': $RANDOM_PASS"
    echo "[init-ssh-config] Set PASSWORD env var to use a custom password"
    echo "============================================"
fi

# Ensure user has a valid shell
usermod -s /bin/bash "$USER_NAME" 2>/dev/null || true

echo "[init-ssh-config] SSH configuration complete"
