# root/etc/s6-overlay/s6-rc.d/svc-ssh/run
#!/command/execlineb -P
# SSH daemon service - runs in foreground for s6 supervision
/usr/sbin/sshd -D -e
